# root/opt/suricata/system/suricata.service
[Unit]
Description=Suricata IDS (Observable Networks distribution)
After=network.target

[Service]
WorkingDirectory=/opt/suricata/
ExecStart=/opt/suricata/run_suricata.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
