# root/opt/suricata/system/suricata.conf
start on (filesystem and net-device-up IFACE!=lo) or (started JOB=network)

chdir /opt/suricata/
respawn

exec /opt/suricata/run_suricata.sh

post-stop exec sleep 10
